# startup.service
[Unit]
Description=LiquidFeedback startup

[Service]
User=root
ExecStart=/opt/lfdeploy/startup.sh

[Install]
WantedBy=multi-user.target